# debian/avahi-dnsconfd.postinst
#!/bin/sh

set -e

case "$1" in
  configure)
    if dpkg --compare-versions "$2" lt-nl "0.8-8"; then
       update-rc.d -f avahi-dnsconfd remove || true
    fi
  ;;
  abort-upgrade|abort-remove|abort-deconfigure)
  ;;
  *)
    echo "postinst called with unknown argument \`$1'" >&2
    exit 1
  ;;
esac

#DEBHELPER#

